# systemSoundDeamon
#!/sbin/openrc-run

name=systemSound
command="/usr/bin/systemSound"
command_args="--my-daemon-args"
command_user="username"
start_stop_daemon_args="--args-for-start-stop-daemon"
command_background="yes"

depend() {
	need alsa
}

start() {
    date > /var/log/systemSound.log
    ebegin "Starting custom deamon (systemSound)"
    start-stop-daemon --start --user username --name systemSound --exec /usr/bin/systemSound startup -- --daemon
    eend $?
}
stop() {
    ebegin "Stoping custom deamon (systemSound)"
    sudo -u username /usr/bin/systemSound shutdown
    eend $?
}
